# ops/systemd/chefflow-prod.service
[Unit]
Description=ChefFlow Production App (app.cheflowhq.com)
Documentation=https://github.com/davidferra13/CFV-1
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=davidferra
Group=davidferra
WorkingDirectory=/home/davidferra/apps/chefflow-prod
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NODE_OPTIONS=--max-old-space-size=1536
EnvironmentFile=/home/davidferra/apps/chefflow-prod/.env.local

# Start command
ExecStart=/home/davidferra/.nvm/versions/node/v20.18.0/bin/node node_modules/.bin/next start -p 3000

# Restart policy
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Memory limits
MemoryMax=2G
MemoryHigh=1800M

# Process settings
StandardOutput=append:/home/davidferra/.pm2/logs/chefflow-prod-out.log
StandardError=append:/home/davidferra/.pm2/logs/chefflow-prod-error.log
SyslogIdentifier=chefflow-prod

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/davidferra/apps/chefflow-prod /home/davidferra/.pm2

[Install]
WantedBy=multi-user.target
