# ops/systemd/chefflow-prod.service
[Unit]
Description=ChefFlow Production App
Documentation=https://github.com/davidferra13/CFV-1
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chefflow
Group=chefflow
WorkingDirectory=/srv/chefflow/current

Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3300
Environment=NODE_OPTIONS=--max-old-space-size=4096
EnvironmentFile=/srv/chefflow/shared/chefflow-prod.env

ExecStart=/usr/bin/env node node_modules/next/dist/bin/next start -p 3300 -H 127.0.0.1

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

MemoryHigh=6G
MemoryMax=8G
TimeoutStartSec=120
TimeoutStopSec=30
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=chefflow-prod

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/srv/chefflow

[Install]
WantedBy=multi-user.target
